# deploy/systemd/aprendaaprogramar-puma.service
[Unit]
Description=aprendaaprogramar Puma (Rails app server)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/aprendaaprogramar/current
Environment=RAILS_ENV=production
EnvironmentFile=-/var/www/aprendaaprogramar/shared/.env
ExecStart=/usr/local/rbenv/shims/bundle exec puma -C /var/www/aprendaaprogramar/current/config/puma.rb
ExecReload=/bin/kill -USR2 $MAINPID
Restart=on-failure
RestartSec=3
SyslogIdentifier=aprendaaprogramar-puma

[Install]
WantedBy=multi-user.target
